# scripts/lifting.service
[Unit]
Description=Lifting Tracker App
After=network.target

[Service]
Type=simple
User=brad
WorkingDirectory=/opt/lifting
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/node packages/server/dist/index.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/lifting

[Install]
WantedBy=multi-user.target
